# roles/mailcatcher/templates/mailcatcher.service.j2
[Unit]
Description=Mailcatcher Service

[Service]
ExecStart=/bin/sh -ec "exec /usr/bin/env $(which mailcatcher) -f -i 0.0.0.0"
Restart=on-failure

[Install]
WantedBy=multi-user.target
